# supervisor.conf
[program:nginx]
command=python /app/leanEngine_app/wsgi.py
