# supervisor.conf
[program:nginx]
directory=/app/leanEngine_app/
command=gunicorn -b 0.0.0.0:9010 -w 20 -k gevent wsgi:application
